# backend/src/worker/.crontab
* * * * * processCloudfarePurgeQueue
* * * * * processAssetDeletionQueue
* * * * * processPendingBulkMessageDeletions
15 * * * * expireAttachments
0 3 * * * userProcessPendingDeletions
0 4 * * 0 processInactivityDeletions
0 23 * * * calculateRetention
